# root/etc/uci-defaults/99_luci-oscam
#!/bin/sh
# 初始化 UCI 配置
uci -q get oscam.main || {
    uci set oscam.main=oscam
    uci set oscam.main.enabled='0'
    uci set oscam.main.pcsc='1'
    uci commit oscam
}
exit 0
